# sources/systemd/scx.service.in
[Unit]
Description=sched-ext BPF CPU Scheduler (@SCHEDULER@)
After=sysinit.target local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c '/usr/bin/scxctl stop 2>/dev/null || true; /usr/bin/scxctl start --sched @SCHEDULER_NAME@ --mode @MODE@'
ExecStop=/usr/bin/scxctl stop
TimeoutStartSec=30

[Install]
WantedBy=sysinit.target
